# qiniu_ufop/data/supervisor.conf
[program:server]
command=qiniu-ufop runserver
loglevel=info
directory=/data/
stderr_logfile=/var/log/server/log

[program:worker]
command=qiniu-ufop runworker
loglevel=info
directory=/data/
stderr_logfile=/var/log/worker/log