# qiniu_ufop/data/supervisor.conf
[program:server]
command=qiniu-ufop runserver
user=ubuntu
loglevel=info
directory=/data/
stderr_logfile=/var/log/server-error

[program:worker]
command=qiniu-ufop runworker
user=ubuntu
loglevel=info
directory=/data/
stderr_logfile=/var/log/worker-error